# systemd/borak.service
[Unit]
Description=BORAK Chat Service
After=network.target ollama.service

[Service]
Type=simple
User=the_bomb
WorkingDirectory=/opt/ollama-ui
Environment="VLLM_ENABLED=false"
Environment="OLLAMA_URL=http://localhost:11434"
ExecStart=/opt/ollama-ui/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8012
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
